# services/nohang-srv/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
exec /usr/bin/nohang --monitor --config /etc/nohang/nohang.conf
