# services/smartmontools-log/run
#!/usr/bin/execlineb -P
foreground { if -n -t { test -d /var/log/smartmontools } install -d -m 0755 -o s6log -g s6log /var/log/smartmontools }
foreground { chmod -R 755 /var/log/smartmontools }
envfile /etc/s6/sv/smartmontools-log/conf
importas -sCiu DIRECTIVES DIRECTIVES
s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /var/log/smartmontools
